# example-initscripts/asnlookup
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: asnlookup
# REQUIRE: LOGIN NETWORKING
# KEYWORD: nojail shutdown
#
# Add the following lines to /etc/rc.conf to enable asnlookup:
# asnlookup_enable (bool):    Set to "NO" by default.
#                             Set it to "YES" to enable asnlookup.

. /etc/rc.subr

name=asnlookup
rcvar=asnlookup_enable

load_rc_config "${name}"

: ${asnlookup_enable:=NO}

command="/usr/local/bin/asn-lookup.py"
pidfile="/tmp/squid-asnbl.pid"
procname="/usr/local/bin/python3.7"
extra_commands=reload
sig_reload=HUP

start_cmd="daemon -u squid /usr/local/bin/asn-lookup.py"

run_rc_command "$1"

# EOF
