# cron.d/reload-asndb.cron
MAILTO=""
25 01 * * * to-be-defined /opt/squid-asnbl/reload_asndb
